# supervisord.sshd.conf
[program:sshd]
autorestart = true
command = /usr/sbin/sshd -D
